# etc/rc.d-script
#!/bin/sh

### BEGIN INIT INFO
# Provides:          hu.dwim.home
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start SBCL with hu.dwim.home at boot time
# Description:       Start the hu.dwim.home services.
### END INIT INFO

. /opt/hu.dwim.home/workspace/hu.dwim.home/bin/environment.sh
. /opt/hu.dwim.home/workspace/hu.dwim.environment/etc/service-scripts/rc.d-script
